# tools/sync-agent/com.jao.sync-trigger.plist.tmpl
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jao.sync-trigger</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>__REPO_DIR__/tools/sync-agent/trigger-handler.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>__REPO_DIR__/.sync-trigger</string>
        <string>__REPO_DIR__/.pull-trigger</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>REPO_DIR</key>
        <string>__REPO_DIR__</string>
        <key>HOME</key>
        <string>__HOME__</string>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>

    <key>ThrottleInterval</key>
    <integer>2</integer>

    <key>StandardOutPath</key>
    <string>__HOME__/Library/Logs/jao-sync-agent.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>__HOME__/Library/Logs/jao-sync-agent.stderr.log</string>
</dict>
</plist>
